# packages/gateway/botdiscord-gateway-staging.service
[Unit]
Description=botdiscord staging gateway (discord.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/botdiscord-staging/packages/gateway
EnvironmentFile=/home/ubuntu/botdiscord-staging/packages/gateway/.env.staging
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
